# flask.conf
description "flask"
start on stopped rc RUNLEVEL=[2345]
respawn
exec export FLASK_APP=/home/adm101/fast-style-transfer2/fast-style-transfer/server.py && flask run --host=0.0.0.0
